# deploy/companybrain-certbot-renew.service
[Unit]
Description=Renew Company Brain Let's Encrypt certificate
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/company-brain
ExecStart=/bin/bash /opt/company-brain/deploy/renew-certificates.sh
